# crontab/crontab_1m.sh
#!/bin/bash
php crontab.php apalancamiento
sleep 0.5 #sleeps 0.5 seconds or 500 miliseconds
php crontab.php apalancamientoCruzado
